# data/units/display-manager-monitor.service
[Unit]
Description=Display manager monitor
After=display-manager.service
StartLimitBurst=5
StartLimitIntervalSec=5

[Service]
Type=simple
SmackProcessLabel=System
User=graphic_fw
Group=graphic_fw
Capabilities=cap_sys_ptrace=i
SecureBits=keep-caps
ExecStart=/usr/bin/bash -c /usr/bin/enlightenment_mon.sh
Restart=always

[Install]
WantedBy=graphical.target

